# services/com.claude-tts.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-tts.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--python</string>
        <string>3.12</string>
        <string>~/.claude-tts/tts-daemon.py</string>
        <string>--foreground</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>WorkingDirectory</key>
    <string>~</string>

    <key>StandardOutPath</key>
    <string>~/.claude-tts/daemon.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>~/.claude-tts/daemon.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:~/.local/bin</string>
    </dict>
</dict>
</plist>
